# rc.local
#!/bin/sh -e

sleep 1

ifdown wlan1
ifup wlan1

echo "1" > /proc/sys/net/ipv4/ip_dynaddr
echo "1" > /proc/sys/net/ipv4/ip_forward

iptables-restore < /opt/iptables.rules
udhcpd /opt/udhcpd.conf
hostapd -B /opt/hostapd.conf
/etc/init.d/nginx start
squid3 -f /opt/squid.conf


exit 0
